# scripts/rebuild-on-webhook.service
[Unit]
Description=Strapi webhook → Astro rebuild
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
Environment=SITE_PATH=/srv/example-site
Environment=WEBHOOK_PORT=9999
# 建议在 Strapi 后台 webhook 设置里添加 header:
#   X-Strapi-Signature: <与此处 WEBHOOK_SECRET 一致的字符串>
Environment=WEBHOOK_SECRET=changeme
ExecStart=/usr/local/bin/rebuild-on-webhook.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
