# scripts/rebuild-on-webhook.service
[Unit]
Description=Strapi webhook → Astro rebuild
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=deploy
Environment=SITE_PATH=/srv/example-site
Environment=WEBHOOK_PORT=9999
# 必须与 Strapi webhook 的 X-Strapi-Signature header 一致，也是 .env.production 里的 WEBHOOK_SECRET
Environment=WEBHOOK_SECRET=changeme
ExecStart=/usr/local/bin/rebuild-on-webhook.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
